# busybox-1.27.2/_install/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/login
